# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /srv/uwsgi.ini
redirect_stderr=true
