# supervisor-app.conf
[program:app-uwsgi]
command=uwsgi --ini /srv/uwsgi.ini
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
